# bikeos/bosd.service
[Unit]
Description=BikeOS WLAN service
After=syslog.target

[Service]
ExecStart=
